# com.openclaw.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Update this path to where you installed the script -->
        <string>/usr/local/bin/openclaw-watchdog.sh</string>
        <string>--daemon</string>
    </array>

    <!-- Start once at load and keep the daemon alive. -->
    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/openclaw-watchdog-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/openclaw-watchdog-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
